# supervisord.conf
[supervisord]
nodaemon=true

[program:xvfb]
priority=10
command=/usr/bin/Xvfb :1 -screen 0 1024x720x16
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb.err

[program:x11vnc]
priority=20
command=x11vnc -display :1 -xkb
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/x11vnc.log
stderr_logfile=/var/log/x11vnc.err

[program:novnc]
priority=25
command=/noVNC/utils/launch.sh
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/novnc.log
stderr_logfile=/var/log/novnc.err

[program:xterm]
priority=30
command=xterm -display :1 -maximized
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/xterm.log
stderr_logfile=/var/log/xterm.err
